# systemd/user/backup-home.timer
[unit]
Description=Run backup-home script every 1h

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
