# systemd/user/rsync-backup.timer
[Unit]
Description=Daily backup of $HOME to /mnt/backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
