# systemd/user/sync-dotfiles.timer
[unit]
Description=Run backup-home script every 1h

[Timer]
OnBootSec=3min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
